# crontab
*/5 * * * * ../code/backup.sh >> /code/cron.log 2>&1
